# deploy/gaz-kotel.service
[Unit]
Description=Gaz Kotel Website (Node.js + PostgreSQL)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/gaz-kotel-website
EnvironmentFile=/var/www/gaz-kotel-website/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /var/www/gaz-kotel-website/server.js
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/gaz-kotel-website
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
